# deploy/linode/gateway-native.service
[Unit]
Description=Native Gateway (Prompt Firewall Benchmark Baseline)
After=network.target

[Service]
Type=simple
ExecStart=/opt/gateway-native/gateway-native
Environment=INFERENCE_URL=http://127.0.0.1:8000
Environment=GATEWAY_REGION=us-ord
Environment=PORT=3000
Environment=RUST_LOG=info
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
